# ops/systemd/classhub-retention.service
[Unit]
Description=ClassHub retention and orphan cleanup
Wants=network-online.target
After=network-online.target docker.service

[Service]
Type=oneshot
# Run as a non-root account that can execute docker compose on this host.
# Override User=/Group= in the deployed unit copy when host identity differs.
User=lms
Group=docker
WorkingDirectory=/srv/lms/app
EnvironmentFile=/srv/lms/app/compose/.env
ExecStartPre=/usr/bin/env bash -lc 'if [[ "$(id -u)" -eq 0 && "${CLASSHUB_ALLOW_ROOT_MAINTENANCE:-0}" != "1" ]]; then echo "Refusing to run classhub-retention as root. Set User=/Group= in the unit (preferred) or set CLASSHUB_ALLOW_ROOT_MAINTENANCE=1 for temporary break-glass."; exit 1; fi'
ExecStart=/usr/bin/env bash /srv/lms/app/scripts/retention_maintenance.sh --compose-mode prod

# Baseline systemd hardening for maintenance scripts.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ProtectControlGroups=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictSUIDSGID=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
ReadWritePaths=/srv/lms/app /run/docker.sock /var/run/docker.sock

# Lower priority so live traffic wins under contention.
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7

[Install]
WantedBy=multi-user.target
